# tools/docker/alignak_container/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:alignak]
command=/bin/sh -c "service alignak start"

[program:nrpe]
command=/bin/sh -c "service nagios-nrpe-server start"
