# tools/docker/shinken_container/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:shinken]
command=/bin/sh -c "service shinken start"
