# tools/docker/surveil_container/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:surveil]
command=/bin/sh -c "surveil-api"

[program:surveil-init]
command=/bin/sh -c "sleep 10 && surveil-init"
